# cron.cron
00 1 * * 5 [ ~/.lastupdate -nt ~/dataPredict.csv ] || ( python3 clasifieur.py ; touch ~/.lastupdate )